# deploy/sub2proxy.service
# Reference unit. deploy/install.sh generates this with the actual PREFIX paths;
# use this file only for a manual install (adjust paths/user to match).
#
# Kept minimal on purpose: filesystem hardening (ProtectSystem/ProtectHome/
# ReadWritePaths) easily makes the data dir read-only unless every path lines up
# with data_dir in config.yaml — a common footgun. NoNewPrivileges is the safe,
# portable bit to keep.
[Unit]
Description=sub2proxy - subscription to multi-port proxy
Documentation=https://github.com/MMortise/Sub2Proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sub2proxy
Group=sub2proxy
WorkingDirectory=/opt/sub2proxy
ExecStart=/opt/sub2proxy/sub2proxy -config /opt/sub2proxy/data/config.yaml
Restart=on-failure
RestartSec=3
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
